# conf/systemd.timer
[Unit]
Description=Run __APP__ regularly

[Timer]
OnCalendar=*:0/8

[Install]
WantedBy=timers.target
